# MONITORING/systemd/node-exporter.service
[Unit]
Description=Prometheus Node Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=monitoring
Group=monitoring
ExecStart=/opt/monitoring/bin/node_exporter \
  --web.listen-address=0.0.0.0:9100 \
  --path.rootfs=/
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
